# src/gnupg/com.dotfiles.gpg-warm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Dotfiles-managed gpg-agent cache warming (see src/gnupg/gpg-warm-agent).

  Runs at login and every 6 hours: preset passphrases expire at max-cache-ttl
  (24h), and GPG Suite's org.gpgtools.macgpg2.shutdown-gpg-agent kills the
  agent at logout, so RunAtLoad re-warms the fresh agent each session. Keeps
  GPG signing TTY-less for agent sessions, hooks, and SSH.

  post_deploy.sh (re)loads this agent on every deploy.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.dotfiles.gpg-warm</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>/Users/{{username}}</string>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/{{username}}/.local/bin/gpg-warm-agent</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>21600</integer>
	<key>StandardOutPath</key>
	<string>/Users/{{username}}/Library/Logs/gpg-warm.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/{{username}}/Library/Logs/gpg-warm.log</string>
</dict>
</plist>
